# dockerfs/etc/services.d/node/run
#!/usr/bin/with-contenv sh

sleep 5
echo "Starting app"
exec node /app/bin/www --dht
