# dockerfs/etc/services.d/lirc/run
#!/usr/bin/with-contenv sh

echo "Starting lircd"
exec lircd -n --driver=default -d /dev/lirc0 -P /dev/null
