# dockerfs/etc/services.d/lirc/run
#!/usr/bin/with-contenv bash

echo "Starting lircd"
exec lircd -n --driver=default -d /dev/lirc0
